# src/timecapsulesmb/assets/boot/samba4/rc.local
#!/bin/sh
set -eu

PATH=/bin:/sbin:/usr/bin:/usr/sbin

/mnt/Flash/start-samba.sh >/dev/null 2>&1 &

if /usr/bin/pkill -0 -f /mnt/Flash/watchdog.sh >/dev/null 2>&1; then
    :
else
    /mnt/Flash/watchdog.sh >/dev/null 2>&1 &
fi

exit 0
