# boot/samba4/rc.local
#!/bin/sh
set -eu

PATH=/bin:/sbin:/usr/bin:/usr/sbin

/mnt/Flash/start-samba.sh >/dev/null 2>&1 &

exit 0
